# crontab
0 * * * * /home/ranjitha/gcp-linux-admin-scripts/check_disk_memory.sh
